# bin/nexus_core.service
[Unit]
Description=NEXUS Cognitive OS Core v3.0
After=network.target

[Service]
Type=simple
User=nexus
WorkingDirectory=/home/zeus/Documentos/ZEUS_SYSTEM
ExecStart=/home/zeus/Documentos/ZEUS_SYSTEM/.venv/bin/python3 -m nexus_core.nexus_core_v3
Restart=always
RestartSec=5
StandardOutput=append:/home/zeus/Documentos/ZEUS_SYSTEM/nexus_core.log
StandardError=append:/home/zeus/Documentos/ZEUS_SYSTEM/nexus_core.log

# Resource Limits (Safety first)
CPUQuota=30%
MemoryLimit=1G

[Install]
WantedBy=multi-user.target
