# deploy/systemd/nexus-root-daemon.service
[Unit]
Description=NEXUS RootDaemon - Secure Command Executor
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/nexus/Documentos/NEXUS/NEXUS
ExecStartPre=/usr/bin/mkdir -p /tmp/nexus
ExecStart=/home/nexus/Documentos/NEXUS/NEXUS/.venv/bin/python3 -m nexus_core.security.root_daemon
Restart=always
RestartSec=3
StandardOutput=append:/home/nexus/Documentos/NEXUS/NEXUS/logs/root_daemon.log
StandardError=append:/home/nexus/Documentos/NEXUS/NEXUS/logs/root_daemon.log

[Install]
WantedBy=default.target
