# deploy/systemd/nexus-cognition.service
[Unit]
Description=NEXUS Cognitive Loop Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/nexus/Documentos/NEXUS/NEXUS
ExecStart=/home/nexus/Documentos/NEXUS/NEXUS/.venv/bin/python -m nexus_core.cognitive.cognition_service
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/nexus/Documentos/NEXUS/NEXUS/.env

[Install]
WantedBy=default.target
